# crontab
* * * * * cronlock python3 /home/app.py >/proc/1/fd/1 2>/proc/1/fd/2


